# archlinux-studio-tweaks/rtc_timer_tweak.service
[Unit]
Description=RTC Timer Tweak

[Service]
ExecStart=/usr/bin/max_rtc_timer.sh

[Install]
WantedBy=multi-user.target
