# grub-btrfs-autogenerate/grub-btrfs-autogenerate.timer
[Unit]
Description=Generate GRUB BTRFS config hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
